# com.music-collector.apple-music-manual.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.music-collector.apple-music-manual</string>

    <key>ProgramArguments</key>
    <array>
        <string>/path/to/music-collector/sync-apple-music.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/path/to/music-collector</string>

    <key>StandardOutPath</key>
    <string>/path/to/music-collector/data/apple_music_manual.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/music-collector/data/apple_music_manual.log</string>

    <key>TimeOut</key>
    <integer>7200</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/home/USERNAME/.local/bin</string>
        <key>PYTHONPATH</key>
        <string>/path/to/music-collector/src</string>
    </dict>
</dict>
</plist>
